# tests/data/test2.cron

# This file was autogenerated by a script. Please do not edit it manually.

30 14 * * * root ./run --debug 
30 6 * * * root ./manage.py restart
