# tests/data/test1.cron
# This file was autogenerated by a script. Please do not edit it manually.
* * * * * root /usr/bin/touch /tmp/test.txt
