# tests/data/test3.cron

# This file was autogenerated by a script. Please do not edit it manually.

MAILTO ='bepetersn@gmailcom'
SHELL  = "/bin/bash"
FAVORITE_ICE_CREAM= $MINT_CHIP

# For more information on the crontab format, 
# see the manual pages of crontab(5) and cron(8).
# m h  dom mon dow
30 14 * * * root ./execute --command-id=1
30 6 * * * root ./execute --command-id=2
